# ansible/templates/nginx-web.service.j2
[Unit]
Description=Nginx web container (CP2)
After=network.target

[Service]
ExecStart=/usr/bin/podman start nginx-web
ExecStop=/usr/bin/podman stop -t 10 nginx-web
Restart=always
RestartSec=10
TimeoutStopSec=70

[Install]
WantedBy=multi-user.target